# tg-max-bridge.service
[Unit]
Description=Telegram to MAX bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/tg-max-bridge
EnvironmentFile=/opt/tg-max-bridge/.env
ExecStart=/opt/tg-max-bridge/venv/bin/python /opt/tg-max-bridge/app.py
Restart=always
RestartSec=5

# Set your Linux user here
User=YOUR_LINUX_USER
Group=YOUR_LINUX_USER

[Install]
WantedBy=multi-user.target
